# systemd/tomato-update-repository.timer
[Unit]
Description=Weekly tomato repository refresh

[Timer]
AccuracySec=6h
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
